# sn-infra-zabbix/sn-zabbix-01-db/backup/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

30 2 * * * /backup/mysqldump_backup.sh >> /var/log/backup.log 2>&1
